# systemd/mcp-guard.timer
[Unit]
Description=Run mcp-guard every 15 minutes

[Timer]
OnCalendar=*:0/15
Persistent=true
AccuracySec=1min

[Install]
WantedBy=timers.target
