# systemd/mcp-guard.timer
[Unit]
Description=Run mcp-guard every 5 minutes

[Timer]
OnCalendar=*:0/5
Persistent=true
AccuracySec=30s

[Install]
WantedBy=timers.target
